# docker/core/root/etc/services.d/cmsLogService/run
#!/usr/bin/with-contenv bash

exec cmsLogService 0
