# docker/core/root/etc/services.d/cmsScoringService/run
#!/usr/bin/with-contenv bash

exec cmsScoringService