# docker/core/root/etc/services.d/cmsResourceService/run
#!/usr/bin/with-contenv bash

exec cmsResourceService 0
