# docker/core/root/etc/services.d/cmsAdminWebServer/run
#!/usr/bin/with-contenv bash

exec cmsAdminWebServer