# docker/rws/root/etc/services.d/cmsRankingWebServer/run
#!/usr/bin/with-contenv bash

exec cmsRankingWebServer
