# docker/core/root/etc/services.d/cmsChecker/run
#!/usr/bin/with-contenv bash

exec cmsChecker