# docker/core/root/etc/services.d/cmsEvaluationService/run
#!/usr/bin/with-contenv bash

contest_id=${CMS_CONTEST_ID}

if [ -n "$contest_id" ]; then
    exec cmsEvaluationService -c "$contest_id" 0
else
    echo | exec cmsEvaluationService 0
fi
