# docker/core/root/etc/services.d/cmsEvaluationService/run
#!/usr/bin/with-contenv bash

contest_id=${CMS_CONTEST_ID}

if [ ! -z "$contest_id" ]; then
    exec cmsEvaluationService -c "$contest_id"
else
    echo | exec cmsEvaluationService
fi
