# docker/worker/root/etc/services.d/cmsWorker/run
#!/usr/bin/with-contenv bash

exec cmsWorker