# docker/worker/root/etc/services.d/cmsWorker/run
#!/usr/bin/with-contenv bash

mkdir -p /run/isolate
echo "/sys/fs/cgroup" > /run/isolate/cgroup
exec cmsWorker
